# OLD/conf/terrastories.conf
description "Start the Terrastories server"
author "Jason Hinebaugh"

start on started docker
stop on stopping docker

exec /bin/terrastories-start
